# roles/pacman/files/pacman-bg.service
[Unit]
Description=Fetch fresh packages in background

[Service]
ExecStart=pacman -Suyw --noconfirm --config /etc/pacman-bg.conf
Restart=on-failure
